# DEBIAN/postinst
# Mise en place du service
systemctl daemon-reload
systemctl enable ha-inet
systemctl start ha-inet
