# DEBIAN/prerm
# Desactivation du service
systemctl stop ha-inet
systemctl disable ha-inet
